# cron/psrm_oracle_watch.cron
*/5 * * * * root /opt/psrm/scripts/check_oracle.sh >/dev/null 2>&1
